# docker/supervisord/tomcat.conf
[supervisord]
nodaemon=true

[program:tomcat]
command=/opt/tomcat/bin/startup.sh
